# packaging/DEBIAN/prerm
#!/bin/sh

systemctl stop nvsm-prometheus-exporters
systemctl disable nvsm-prometheus-exporters
